# ops/systemd/solefeed-campground.slice
[Unit]
Description=Campground Telegram bot slice

[Slice]
CPUQuota=150%
CPUWeight=20
IOWeight=50
MemoryHigh=3072M
MemoryMax=4096M
TasksMax=256
